# infrastructure/chef-repo/cookbooks/thena/templates/default/thena-uwsgi.conf.erb
description "uWSGI server instance configured to serve thena"
author "Daniel Kronovet"

start on runlevel [2345]
stop on runlevel [!2345]

chdir <%= node['thena']['config_dir'] %>
exec uwsgi --ini thena-uwsgi.ini
