# infrastructure/thena.conf
description "uWSGI server instance configured to serve thena"

start on runlevel [2345]
stop on runlevel [!2345]

setuid ubuntu
setgid www-data

chdir /srv/www/thena/thena
exec uwsgi --ini thena.ini